# systemctl/orch-celery-fileapp-worker.service
[Unit]
Description=Orch Celery Worker (FileApp Ingest/Process)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=orch
Group=orch
WorkingDirectory=/opt/orch

EnvironmentFile=/etc/orch/orch.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/orch
Environment=ORCH_QUEUE_PROFILE=prod
Environment=CELERY_ENABLED=true
Environment=CELERY_FILEAPP_INGEST_ENABLED=true
Environment=CELERY_GENERATE_FILE_ENABLED=false
Environment=CELERY_S3_FILES_INGEST_QUEUE=orch_fileapp_ingest_events
Environment=CELERY_SOURCE_LIST_INGEST_QUEUE=orch_fileapp_source_list_ingest

ExecStart=/opt/orch/.venv/bin/celery -A app.core.celery_app:celery_app worker \
  --hostname=orch-celery-fileapp-worker@136_01 \
  --loglevel=INFO \
  --pool=prefork \
  --queues=orch_fileapp_ingest_events,orch_fileapp_source_list_ingest \
  --concurrency=4 \
  --prefetch-multiplier=1 \
  -Ofair \
  --max-tasks-per-child=1000 \
  --max-memory-per-child=350000 \
  --without-gossip \
  --without-mingle

Restart=always
RestartSec=5
TimeoutStopSec=180
KillMode=mixed

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
