# launchd/com.orch.celery.beat.generate_file.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.orch.celery.beat.generate_file</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd "__REPO_ROOT__" &amp;&amp; source .venv/bin/activate &amp;&amp; export CELERY_ENABLED=true WORKFLOW_V2_ENABLED=true WORKFLOW_V2_EXECUTE_M2=true CELERY_BEAT_DISPATCH_ENABLED=false CELERY_BEAT_HEARTBEAT_ENABLED=false CELERY_GENERATE_FILE_ENABLED=true CELERY_GENERATE_FILE_SCAN_ENABLED=true CELERY_GENERATE_FILE_SCAN_INTERVAL_SECONDS=10 CELERY_GENERATE_FILE_WORKSPACE_UUID=ba7eb0ec-e565-447c-8c11-8f870cf72a60 CELERY_GENERATE_FILE_RUN_QUEUE=orch_component_generate_file_run_launchd_local CELERY_GENERATE_FILE_SCAN_QUEUE=orch_component_generate_file_scan_launchd_local &amp;&amp; exec celery -A app.core.celery_app:celery_app beat --schedule=/tmp/orch-celerybeat-gf-launchd-local -l INFO</string>
  </array>
  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/.runlogs/launchd/beat_gf.out.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/.runlogs/launchd/beat_gf.err.log</string>
</dict>
</plist>
